# relay-module.service
[Unit]
Description=Raspberry Pi Relay Module REST API Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=jodok
Group=jodok
WorkingDirectory=/home/jodok/sandbox/raspberry-pi-relay-module
ExecStart=/home/jodok/sandbox/raspberry-pi-relay-module/venv/bin/python /home/jodok/sandbox/raspberry-pi-relay-module/relay_module.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raspberry-pi-relay

# GPIO access requires special permissions
# The install script will set up udev rules for GPIO access

[Install]
WantedBy=multi-user.target
